# docs/crontab
#crontab -u www -e
#定时清空商户及渠道单日限额
01 0 * * *  /www/wwwroot/gateway.payment.com/protected/yii account/reset-quota  2>&1